# templates/stash.service.erb
[Unit]
Description=Atlassian Systemd Jira Service
After=syslog.target network.target

[Service]
Type=forking
Environment="JAVA_HOME=<%= scope.lookupvar('stash::javahome') %>"
<%- if scope.lookupvar('extra_environment') -%>
<%- scope.lookupvar('extra_environment').each do |key, val| -%>
Environment="<%= key %>=<%= val %>"
<% end -%>
<% end -%>
UMask=0027
PIDFile=<%= scope.lookupvar('pidfile') %>
User=<%= scope.lookupvar('stash::user') %>
ExecStart=<%= scope.lookupvar('stash::webappdir') %>/bin/start-<%= scope.lookupvar('stash::product') %>.sh
ExecStop=<%= scope.lookupvar('stash::webappdir') %>/bin/stop-<%= scope.lookupvar('stash::product') %>.sh

[Install]
WantedBy=multi-user.target
